# backend/crontab
*/1 * * * * root cd /app && python manage.py import >> /var/log/cron.log 2>&1
